# build/src/supervisord.conf
[supervisord]
nodaemon=true
#user=dhcp
logfile = /data/supervisord.log
logfile_maxbytes = 0
logfile_backups=0
loglevel = info
pidfile = /data/supervisord.pid

[program:dhcp_server]
command=/usr/sbin/dnsmasq --conf-file=/data/conf/dnsmasq.conf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/data/supervisord.log
stdout_logfile_maxbytes=0

